# root/etc/services.d/apache/run
#!/usr/bin/with-contenv bash

[ "$MODE" == "MASTER" ] && exec /usr/sbin/apachectl -D FOREGROUND
